# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:airflow]
command=/bin/bash -c "/root/run_php_fpm.sh"
